# pochtahub.service
[Unit]
Description=PochtaHub Django Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/pochtahub
Environment="PATH=/var/www/pochtahub/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/bin/bash -c 'cd /var/www/pochtahub && source venv/bin/activate && python manage.py runserver 0.0.0.0:8000'
Restart=always
RestartSec=10
StandardOutput=append:/var/www/pochtahub/server.log
StandardError=append:/var/www/pochtahub/server.log

[Install]
WantedBy=multi-user.target

